# paas.conf
start on startup
stop on shutdown
respawn
respawn limit 10 5

script
	export NODE_ENV="production" 
	sudo -u cozy haibu --coffee
end script