# paas.conf
start on startup
stop on shutdown
respawn
respawn limit 10 5

script
NODE_ENV="production" haibu --coffee
end script