# paas.conf
start on startup
stop on shutdown
respawn
respawn limit 10 5

script
	export NODE_ENV="production" 
	su -c '/home/cozy/cozy-setup/node_modules/haibu/bin/haibu --coffee' cozy
end script